# site-cookbooks/fort-pinecone/templates/default/rocketchat.service.erb
[Unit]
Description=The Rocket.Chat server
After=network.target remote-fs.target nss-lookup.target nginx.target mongod.target
[Service]
ExecStart=/usr/local/bin/node /opt/Rocket.Chat/main.js
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rocketchat
User=rocketchat
Environment=MONGO_URL=<%= @mongo %> ROOT_URL=<%= @root %> PORT=<%= @port %>
[Install]
WantedBy=multi-user.target